# deploy/systemd/qiaolian-user-bot.service
[Unit]
Description=Qiaolian User Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qiaolianbot
Group=qiaolianbot
WorkingDirectory=/opt/qiaolian_dual_bots
EnvironmentFile=/opt/qiaolian_dual_bots/.env
ExecStart=/opt/qiaolian_dual_bots/.venv/bin/python /opt/qiaolian_dual_bots/run_user_bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
