# deploy/systemd/qiaolian-admin-web.service
[Unit]
Description=Qiaolian Admin Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qiaolianbot
Group=qiaolianbot
WorkingDirectory=/opt/qiaolian_dual_bots/v2_admin
EnvironmentFile=/opt/qiaolian_dual_bots/.env
ExecStart=/opt/qiaolian_dual_bots/.venv/bin/waitress-serve --host=0.0.0.0 --port=${ADMIN_PORT} admin_server:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
