# deploy/systemd/qiaolian-publisher-bot.service
[Unit]
Description=Qiaolian Publisher Bot v2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qiaolianbot
Group=qiaolianbot
WorkingDirectory=/opt/qiaolian_dual_bots
EnvironmentFile=/opt/qiaolian_dual_bots/.env
ExecStart=/opt/qiaolian_dual_bots/.venv/bin/python /opt/qiaolian_dual_bots/v2/run_publisher_bot_v2.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
